# .docker/service/ping/run
#!/bin/bash
echo "Running ping runner."
while true; do
   /app/ping.php
   sleep 5s
done