# .docker/service/speedtest/run
#!/bin/bash
echo "Running speedtest runner."
while true; do
   /app/speedtest.sh
   sleep 5m
done